# base_external/rootfs_overlay/etc/init.d/S98lddmodules
#!/bin/sh
### BEGIN INIT INFO
# Provides:          custom_modules
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Script to load and unload multiple modules
# Description:       This script will load and unload multiple modules
### END INIT INFO

# Define the names of the modules
MODULES="scull hello faulty"

# Check if a module is already loaded
is_loaded() {
    module=$1
    lsmod | grep "^$module " >/dev/null 2>&1
}

# Create device nodes (excluding hello module)
create_device_nodes() {
    module=$1
    if [ "$module" != "hello" ]; then
        device=$module
        major=$(awk "\$2==\"$module\" {print \$1}" /proc/devices)
        if [ -z "$major" ]; then
            echo "Failed to get major number for $module"
            exit 1
        fi
        if [ ! -e /dev/$device ]; then
            echo "Creating /dev/$device with major number $major"
            mknod /dev/$device c $major 0
        else
            echo "/dev/$device already exists"
        fi
    fi
}

# Load the modules
load_modules() {
    for module in $MODULES; do
        if is_loaded "$module"; then
            echo "$module is already loaded"
        else
            echo "Loading $module module"
            modprobe "$module"
            if [ $? -ne 0 ]; then
                echo "Failed to load $module module"
                exit 1
            fi
            create_device_nodes "$module"
        fi
    done
}

# Unload the modules and remove device nodes
unload_modules() {
    for module in $MODULES; do
        if is_loaded "$module"; then
            echo "Unloading $module module"
            modprobe -r "$module"
            if [ $? -ne 0 ]; then
                echo "Failed to unload $module module"
                exit 1
            fi
            if [ "$module" != "hello" ]; then
                echo "Removing /dev/$module"
                rm -f /dev/$module
            fi
        else
            echo "$module is not loaded"
        fi
    done
}

# Parse the command line arguments
case "$1" in
    start)
        load_modules
        ;;
    stop)
        unload_modules
        ;;
    restart)
        unload_modules
        load_modules
        ;;
    status)
        for module in $MODULES; do
            if is_loaded "$module"; then
                echo "$module is loaded"
            else
                echo "$module is not loaded"
            fi
        done
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0

